# .config/systemd/user/battery-notify.service
[Unit]
Description=Battery threshold notifications (80%% charge / 20%% low)

[Service]
Type=oneshot
ExecStart=%h/.config/scripts/battery-notify.sh
